# build/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:webservice]
command=java -jar HelloWorldWebService-0.0.1-SNAPSHOT.jar
directorv=/WebService
autostart=true
autorestart=true
stderr_logfile=/var/log/WebService/service.err.log
stdout_logfile=/var/log/WebService/service.out.log